# deploy/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
logfile_backups=0
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
silent=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:comfyui]
command=python main.py --listen %(ENV_COMFYUI_SERVER_HOST)s --port %(ENV_COMFYUI_SERVER_PORT)s
directory=/root/ComfyUI
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=PYTHONUNBUFFERED=1
priority=100
startsecs=10
startretries=3

[program:handler]
command=python server.py
directory=/app/comfyui-sync-handler
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=PYTHONUNBUFFERED=1,COMFYUI_SERVER_HOST=%(ENV_COMFYUI_SERVER_HOST)s,COMFYUI_SERVER_PORT=%(ENV_COMFYUI_SERVER_PORT)s,SERVER_HOST=%(ENV_SERVER_HOST)s,SERVER_PORT=%(ENV_SERVER_PORT)s,DEFAULT_TIMEOUT=%(ENV_DEFAULT_TIMEOUT)s,LOG_LEVEL=%(ENV_LOG_LEVEL)s
priority=200
startsecs=15
startretries=3
depends_on=comfyui

[group:comfyui-services]
programs=comfyui,handler
priority=999 